# modbus/rootfs/etc/services.d/modbus/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start modbus service if enabled
# ==============================================================================

USBPORT=$(bashio::config 'usbport')
MQTTHOST=$(bashio::config 'mqtthost')
MQTTUSER=$(bashio::config 'mqttuser')
MQTTPASSWORD=$(bashio::config 'mqttpassword')

bashio::log.info "Starting the Modbus daemon..."
bashio::log.info "USB port is.........: ${USBPORT}"
bashio::log.info "MQTT host is........: ${MQTTHOST}"
bashio::log.info "MQTT user is........: ${MQTTUSER}"
bashio::log.info "MQTT password is....: ${MQTTPASSWORD}"
exec /usr/bin/poll_isystem_mqtt.py --serial ${USBPORT} --user ${MQTTUSER} --password ${MQTTPASSWORD} ${MQTTHOST}
